# server/opal-api.service.example
[Unit]
Description=opal-api backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/opal-api
EnvironmentFile=/opt/opal-api/.env
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
User=opal

[Install]
WantedBy=multi-user.target
